# supervisord/supervisor-app.conf
[supervisord]
nodaemon=false

[program:apache2]
command=/usr/local/bin/httpd-foreground

[program:shibd]
command = service shibd start